# web/config/sv/supervisor.conf
[supervisord]
nodaemon=true

[program:cron]
command=cron -f

[program:service]
command=bash -c "python manage.py collectstatic --no-input && gunicorn project.wsgi --bind :8000"
